# bin/claude-otel-stack.service
# Generated By: Claude Code (Claude Opus 4.6)
# systemd user unit for claude-otel-stack
#
# Install:
#   cp claude-otel-stack.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable claude-otel-stack
#   systemctl --user start claude-otel-stack
#
# Logs:
#   journalctl --user -u claude-otel-stack -f

[Unit]
Description=Claude Code OTEL monitoring stack
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/source/claude-otel-stack
ExecStart=%h/.local/bin/podman-compose -f docker-compose.yml up
ExecStop=%h/.local/bin/podman-compose -f docker-compose.yml down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
