# bin/otel-claude.service
# Generated By: Claude Code (Claude Opus 4.6)
# systemd user unit for otel-claude stack
#
# Install:
#   cp otel-claude.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable otel-claude
#   systemctl --user start otel-claude
#
# Logs:
#   journalctl --user -u otel-claude -f

[Unit]
Description=OTEL Claude Code monitoring stack
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/source/claude-otel-stack
ExecStart=/usr/bin/podman-compose -f docker-compose.yml up
ExecStop=/usr/bin/podman-compose -f docker-compose.yml down
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
